# deploy/gated-deploy.service
[Unit]
Description=Deploy latest GATED release from GitHub
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment="GITHUB_REPO=KernPhilipp/GATED"
EnvironmentFile=-/home/gated/gated-frontend/shared/deploy.env
ExecStart=/home/gated/gated-frontend/bin/gated-deploy.sh
