# deploy/gated-backend.service
[Unit]
Description=GATED Dart backend
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gated
Group=gated
WorkingDirectory=/home/gated/gated-frontend/current/backend
EnvironmentFile=/home/gated/gated-frontend/shared/backend.env
ExecStart=/usr/bin/env dart run server.dart
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
